# scripts/autoupdate/autoupdate.service
[Unit]
Description=Pacman Autoupdate

[Service]
ExecStart=/home/x/scripts/autoupdate/autoupdate

[Install]
WantedBy=multi-user.target
