# scripts/autoupdate/autoupdate.service
[Unit]
Description=Pacman Autoupdate

[Service]
ExecStart=/bin/sh /home/x/scripts/autoupdate.sh

[Install]
WantedBy=multi-user.target
